# status-mail.timer
[Unit]
Description=Regularly check system status

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1m

[Install]
WantedBy=timers.target
